# debian/shell/backup/systemd/restic_backup.service
[Unit]
Description=Restic backup service

[Service]
Type=oneshot
ExecStart=/root/restic_backup/restic_backup_helper.sh